# hausofsutra/crontab
# Run every 12 hours at midnight and noon
# Source environment variables that the entrypoint exported for cron
0 0,12 * * * . /app/sync-env.sh && /app/sync.sh >> /var/log/hausofsutra/sync.log 2>&1
